# debian/davis-wx-web.service
[Unit]
Description=Davis Weather Station Web UI
After=network.target davis-wx-logger.service
Wants=davis-wx-logger.service

[Service]
Type=simple
User=davis
Group=davis
WorkingDirectory=/opt/davis-wx/backend
EnvironmentFile=/etc/davis-wx/davis-wx.conf
ExecStart=/opt/davis-wx/backend/.venv/bin/python -m uvicorn app.main:app --host ${DAVIS_HOST} --port ${DAVIS_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
