# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create davis system user if it doesn't exist
        if ! id -u davis >/dev/null 2>&1; then
            adduser --system --group --home /var/lib/davis-wx \
                    --no-create-home --shell /usr/sbin/nologin davis
        fi

        # Add to dialout for serial port access
        usermod -aG dialout davis 2>/dev/null || true

        # Create Python venv and install dependencies
        VENV=/opt/davis-wx/backend/.venv
        echo "Creating Python virtual environment..."
        python3 -m venv "$VENV"
        echo "Installing Python dependencies (this may take a minute)..."
        "$VENV/bin/pip" install --upgrade pip --quiet 2>/dev/null
        "$VENV/bin/pip" install /opt/davis-wx/backend/ --quiet

        # Fix ownership
        chown -R davis:davis /var/lib/davis-wx
        chown -R davis:davis /opt/davis-wx

        # Enable and start services
        systemctl daemon-reload
        systemctl enable davis-wx-logger.service davis-wx-web.service
        systemctl start davis-wx-logger.service davis-wx-web.service || true
        ;;
esac

#DEBHELPER#
exit 0
